# data/systemd/dactl-daq.service.in
[Unit]
Description=Data Acquisition Daemon
After=network.target
After=dbus.service
Conflicts=
Documentation=man:dactl-daq(1)

[Service]
EnvironmentFile=-@sysconfdir@/dactl
ExecStart=@sbindir@/dactl-daq $DAQD_ARGS
ExecReload=/bin/kill -HUP $MAINPID
# supress to log debug and error also to /var/log/messages
StandardOutput=null
StandardError=null
Type=dbus
BusName=org.coanda.Dactl.DAQ

[Install]
WantedBy=basic.target
Alias=dbus-org.coanda.Dactl.DAQ.service
