# data/systemd/dactl-ctl.service.in
[Unit]
Description=Plant Control Daemon
After=network.target
After=dbus.service
Conflicts=
Documentation=man:dactl-ctl(1)

[Service]
EnvironmentFile=-@sysconfdir@/dactl
ExecStart=@sbindir@/dactl-ctl $CTLD_ARGS
ExecReload=/bin/kill -HUP $MAINPID
# supress to log debug and error also to /var/log/messages
StandardOutput=null
StandardError=null
Type=dbus
BusName=org.coanda.Dactl.Control

[Install]
WantedBy=basic.target
Alias=dbus-org.coanda.Dactl.Control.service
